# security/fail2ban/files/fail2ban.rc
#!/bin/sh

daemon="${TRUEPREFIX}/bin/fail2ban-client"
daemon_flags="start"
#daemon_user="root"

. /etc/rc.d/rc.subr
